# systemd/360blur.service
[Unit]
Description=360blur Video Processing Service
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%INSTALLDIR%
ExecStart=%INSTALLDIR%/venv/bin/python %INSTALLDIR%/blur360_webapp.py
Restart=on-failure
RestartSec=5s
Environment=PYTHONUNBUFFERED=1

# Security options
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target